# supervisord.conf
[supervisord]
nodaemon=true

[program:nodeapp]
command=node server/index.js
directory=/app
autostart=true
autorestart=true

[program:nginx]
command=nginx -g "daemon off;"